# tools/deployment/package/single-container/supervisord.conf
[supervisord]
nodaemon=true
loglevel=warn
logfile=/var/log/clp-single-container/supervisord.log
pidfile=/run/clp-single-container-supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:mariadb]
command=/usr/sbin/mariadbd --defaults-extra-file=/opt/clp/etc/mysql/conf.d/logging.cnf --user=mysql --bind-address=127.0.0.1 --port=3306 --datadir=/var/data/database
autorestart=true
priority=10
stdout_logfile=/var/log/clp-single-container/mariadb.stdout.log
stderr_logfile=/var/log/clp-single-container/mariadb.stderr.log

[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
autorestart=true
priority=20
stdout_logfile=/var/log/clp-single-container/rabbitmq.stdout.log
stderr_logfile=/var/log/clp-single-container/rabbitmq.stderr.log
user=rabbitmq
environment=HOME="/var/lib/rabbitmq",RABBITMQ_LOGS="/var/log/rabbitmq/rabbitmq.log",ERL_EPMD_ADDRESS="127.0.0.1",RABBITMQ_SERVER_ADDITIONAL_ERL_ARGS="-kernel inet_dist_use_interface {127,0,0,1}"

[program:redis]
command=/usr/bin/redis-server /opt/clp/etc/redis/redis.conf --dir /var/data/redis --protected-mode no --bind 127.0.0.1 --requirepass %(ENV_REDIS_PASSWORD)s
autorestart=true
priority=30
stdout_logfile=/var/log/clp-single-container/redis.stdout.log
stderr_logfile=/var/log/clp-single-container/redis.stderr.log
user=redis

[program:mongodb]
command=/usr/bin/mongod --config /opt/clp/etc/mongo/mongod.conf --dbpath /var/data/results_cache --bind_ip 127.0.0.1
autorestart=true
priority=40
stdout_logfile=/var/log/clp-single-container/mongodb.stdout.log
stderr_logfile=/var/log/clp-single-container/mongodb.stderr.log
user=mongodb

[program:db-table-creator]
command=/bin/bash -c '/usr/local/bin/clp-single-container-wait-for-tcp 127.0.0.1 3306 python3 -u -m clp_py_utils.create-db-tables --config %(ENV_CLP_CONFIG_PATH)s --storage-engine %(ENV_CLP_PACKAGE_STORAGE_ENGINE)s && touch /run/clp-single-container/db-table-creator.done'
autorestart=false
startsecs=0
priority=100
stdout_logfile=/var/log/clp-single-container/db-table-creator.stdout.log
stderr_logfile=/var/log/clp-single-container/db-table-creator.stderr.log
user=clp-user

[program:results-cache-indices-creator]
command=/bin/bash -c '/usr/local/bin/clp-single-container-wait-for-tcp 127.0.0.1 27017 python3 -u -m clp_py_utils.initialize-results-cache --uri mongodb://results_cache:%(ENV_CLP_RESULTS_CACHE_CONNECT_PORT)s/%(ENV_CLP_RESULTS_CACHE_DB_NAME)s --stream-collection %(ENV_CLP_RESULTS_CACHE_STREAM_COLLECTION_NAME)s && touch /run/clp-single-container/results-cache-indices-creator.done'
autorestart=false
startsecs=0
priority=105
stdout_logfile=/var/log/clp-single-container/results-cache-indices-creator.stdout.log
stderr_logfile=/var/log/clp-single-container/results-cache-indices-creator.stderr.log
user=clp-user

[program:compression-scheduler]
command=/bin/bash -c 'export RESULT_BACKEND="redis://default:${CLP_REDIS_PASS:?Please set CLP_REDIS_PASS.}@redis:${CLP_REDIS_CONNECT_PORT}/${CLP_REDIS_BACKEND_DB_COMPRESSION:-1}"; exec /usr/local/bin/clp-single-container-wait-for-files /run/clp-single-container/db-table-creator.done -- /usr/local/bin/clp-single-container-wait-for-tcp 127.0.0.1 5672 /usr/local/bin/clp-single-container-wait-for-tcp 127.0.0.1 6379 python3 -u -m job_orchestration.scheduler.compress.compression_scheduler --config %(ENV_CLP_CONFIG_PATH)s'
autorestart=true
priority=110
stdout_logfile=/var/log/clp-single-container/compression-scheduler.stdout.log
stderr_logfile=/var/log/clp-single-container/compression-scheduler.stderr.log
user=clp-user

[program:compression-worker]
command=/bin/bash -c 'export RESULT_BACKEND="redis://default:${CLP_REDIS_PASS:?Please set CLP_REDIS_PASS.}@redis:${CLP_REDIS_CONNECT_PORT}/${CLP_REDIS_BACKEND_DB_COMPRESSION:-1}"; exec /usr/local/bin/clp-single-container-wait-for-tcp 127.0.0.1 5672 /usr/local/bin/clp-single-container-wait-for-tcp 127.0.0.1 6379 python3 -u /opt/clp/lib/python3/site-packages/bin/celery -A job_orchestration.executor.compress worker --concurrency %(ENV_CLP_COMPRESSION_WORKER_CONCURRENCY)s --loglevel WARNING -f /var/log/compression_worker/worker.log -Q compression -n compression-worker@%%h'
autorestart=true
priority=120
stdout_logfile=/var/log/clp-single-container/compression-worker.stdout.log
stderr_logfile=/var/log/clp-single-container/compression-worker.stderr.log
user=clp-user

[program:query-scheduler]
command=/bin/bash -c 'export RESULT_BACKEND="redis://default:${CLP_REDIS_PASS:?Please set CLP_REDIS_PASS.}@redis:${CLP_REDIS_CONNECT_PORT}/${CLP_REDIS_BACKEND_DB_QUERY:-0}"; exec /usr/local/bin/clp-single-container-wait-for-files /run/clp-single-container/db-table-creator.done -- /usr/local/bin/clp-single-container-wait-for-tcp 127.0.0.1 5672 /usr/local/bin/clp-single-container-wait-for-tcp 127.0.0.1 6379 python3 -u -m job_orchestration.scheduler.query.query_scheduler --config %(ENV_CLP_CONFIG_PATH)s'
autorestart=true
priority=130
stdout_logfile=/var/log/clp-single-container/query-scheduler.stdout.log
stderr_logfile=/var/log/clp-single-container/query-scheduler.stderr.log
user=clp-user

[program:query-worker]
command=/bin/bash -c 'export RESULT_BACKEND="redis://default:${CLP_REDIS_PASS:?Please set CLP_REDIS_PASS.}@redis:${CLP_REDIS_CONNECT_PORT}/${CLP_REDIS_BACKEND_DB_QUERY:-0}"; exec /usr/local/bin/clp-single-container-wait-for-tcp 127.0.0.1 5672 /usr/local/bin/clp-single-container-wait-for-tcp 127.0.0.1 6379 python3 -u /opt/clp/lib/python3/site-packages/bin/celery -A job_orchestration.executor.query worker --concurrency %(ENV_CLP_QUERY_WORKER_CONCURRENCY)s --loglevel WARNING -f /var/log/query_worker/worker.log -Q query -n query-worker@%%h'
autorestart=true
priority=140
stdout_logfile=/var/log/clp-single-container/query-worker.stdout.log
stderr_logfile=/var/log/clp-single-container/query-worker.stderr.log
user=clp-user

[program:reducer]
command=/usr/local/bin/clp-single-container-wait-for-files /run/clp-single-container/results-cache-indices-creator.done -- /usr/local/bin/clp-single-container-wait-for-tcp 127.0.0.1 7000 python3 -u -m job_orchestration.reducer.reducer --config %(ENV_CLP_CONFIG_PATH)s --concurrency %(ENV_CLP_REDUCER_CONCURRENCY)s --upsert-interval %(ENV_CLP_REDUCER_UPSERT_INTERVAL)s
autorestart=true
priority=150
stdout_logfile=/var/log/clp-single-container/reducer.stdout.log
stderr_logfile=/var/log/clp-single-container/reducer.stderr.log
user=clp-user

[program:garbage-collector]
command=/bin/bash -c 'if [[ "${CLP_GARBAGE_COLLECTOR_ENABLED:-1}" == "0" ]]; then echo "garbage-collector disabled."; exit 0; fi; exec /usr/local/bin/clp-single-container-wait-for-files /run/clp-single-container/db-table-creator.done /run/clp-single-container/results-cache-indices-creator.done -- python3 -u -m job_orchestration.garbage_collector.garbage_collector --config %(ENV_CLP_CONFIG_PATH)s'
autorestart=unexpected
startsecs=0
priority=155
stdout_logfile=/var/log/clp-single-container/garbage-collector.stdout.log
stderr_logfile=/var/log/clp-single-container/garbage-collector.stderr.log
user=clp-user
environment=CLP_LOGGING_LEVEL="%(ENV_CLP_GARBAGE_COLLECTOR_LOGGING_LEVEL)s",CLP_LOGS_DIR="/var/log/garbage_collector"

[program:api-server]
command=/bin/bash -c 'if [[ "${CLP_API_SERVER_ENABLED:-1}" == "0" ]]; then echo "api-server disabled."; exit 0; fi; exec /usr/local/bin/clp-single-container-wait-for-files /run/clp-single-container/db-table-creator.done /run/clp-single-container/results-cache-indices-creator.done -- /opt/clp/bin/api_server --host 0.0.0.0 --port 3001 --config %(ENV_CLP_CONFIG_PATH)s'
autorestart=unexpected
startsecs=0
priority=160
stdout_logfile=/var/log/clp-single-container/api-server.stdout.log
stderr_logfile=/var/log/clp-single-container/api-server.stderr.log
user=clp-user

[program:log-ingestor]
command=/bin/bash -c 'if [[ "${CLP_LOG_INGESTOR_ENABLED:-0}" == "0" ]]; then echo "log-ingestor disabled."; exit 0; fi; exec /usr/local/bin/clp-single-container-wait-for-files /run/clp-single-container/db-table-creator.done -- /opt/clp/bin/log-ingestor --config %(ENV_CLP_CONFIG_PATH)s --host 0.0.0.0 --port 3002'
autorestart=unexpected
startsecs=0
priority=165
stdout_logfile=/var/log/clp-single-container/log-ingestor.stdout.log
stderr_logfile=/var/log/clp-single-container/log-ingestor.stderr.log
user=clp-user
environment=RUST_LOG="%(ENV_CLP_LOG_INGESTOR_LOGGING_LEVEL)s"

[program:mcp-server]
command=/bin/bash -c 'if [[ "${CLP_MCP_SERVER_ENABLED:-0}" == "0" ]]; then echo "MCP server disabled."; exit 0; fi; exec /usr/local/bin/clp-single-container-wait-for-files /run/clp-single-container/db-table-creator.done /run/clp-single-container/results-cache-indices-creator.done -- python3 -u -m clp_mcp_server.clp_mcp_server --host 0.0.0.0 --port 8000 --config-path %(ENV_CLP_CONFIG_PATH)s'
autorestart=unexpected
startsecs=0
priority=168
stdout_logfile=/var/log/clp-single-container/mcp-server.stdout.log
stderr_logfile=/var/log/clp-single-container/mcp-server.stderr.log
user=clp-user
environment=CLP_LOGGING_LEVEL="%(ENV_CLP_MCP_LOGGING_LEVEL)s",CLP_LOGS_DIR="/var/log/mcp_server"

[program:webui]
command=/usr/local/bin/clp-single-container-wait-for-files /run/clp-single-container/db-table-creator.done /run/clp-single-container/results-cache-indices-creator.done -- /opt/clp/bin/node-22 /opt/clp/var/www/webui/server/dist/src/main.js
autorestart=true
priority=170
stdout_logfile=/var/log/clp-single-container/webui.stdout.log
stderr_logfile=/var/log/clp-single-container/webui.stderr.log
user=clp-user
